# ops/worker/lyrashield-worker-secrets.service
[Unit]
Description=Load LyraShield worker secrets through the VM managed identity
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=AZURE_CONFIG_DIR=/var/lib/lyrashield-azure
ExecStart=/usr/local/libexec/lyrashield-refresh-secrets
RemainAfterExit=yes
